# cloud-minder.service.template
[Unit]
Description=CloudMinder
After=network.target
[Service]
User=__USER__
WorkingDirectory=/home/__USER__/CloudMinder
ExecStart=__NODE__ CloudMinder
[Install]
WantedBy=multi-user.target

